# daemon/com.triadblue.linksblue-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.triadblue.linksblue-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>/Users/deanlewis/linksblue-network/daemon/daemon.js</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>WorkingDirectory</key>
    <string>/Users/deanlewis/linksblue-network/daemon</string>

    <key>StandardOutPath</key>
    <string>/Users/deanlewis/.linksblue-daemon/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/deanlewis/.linksblue-daemon/daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
